# .crontab
*/1 * * * * php ~/apps/php3/shedule.php >> ~/apps/php3/storage/logs/cron.log

