# .crontab
*/1 * * * * artisan check:visits vasya vasya321@mail.ru >> ~/apps/php3/storage/logs/cron.log

