# .crontab
*/1 * * * * echo "hello world" >> ~/apps/php3/storage/logs/cron.log

